# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:Frozone]
command=/frozone-bin/Frozone config.yml
directory=/frozone-bin
stdout_logfile=/var/log/frozone/stdout
stderr_logfile=/var/log/frozone/stderr
